# systemd/archive-cron-hourly.timer
[Unit]
Description=Archive Hourly Cleanup Timer

[Timer]
OnCalendar=*-*-* *:0/10:00

[Install]
WantedBy=timers.target
